# server/etc/systemd/system/lastcontrol-securitydata-update.service
[Unit]
Description=LastControl SecurityData Update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lastcontrol-securitydata-update.sh
